# desktop/systemd/ups-orchestrator-agent.service
[Unit]
Description=UPS Orchestrator Desktop Agent
After=graphical-session.target network-online.target ups-power-agent.service
Wants=graphical-session.target network-online.target
Requires=ups-power-agent.service
PartOf=graphical-session.target

[Service]
Type=simple
WorkingDirectory=/opt/ups-orchestrator/desktop
EnvironmentFile=/opt/ups-orchestrator/desktop/.env

ExecStart=/opt/ups-orchestrator/desktop/venv/bin/python /opt/ups-orchestrator/desktop/app/agent.py
ExecStop=/opt/ups-orchestrator/desktop/venv/bin/python /opt/ups-orchestrator/desktop/app/send_state.py offline

Restart=always
RestartSec=2

[Install]
WantedBy=graphical-session.target